# print-server/print-server.service
[Unit]
Description=The Breakery Print Server
After=network.target

[Service]
Type=simple
User=printserver
WorkingDirectory=/opt/breakery/print-server
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=print-server

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
